# roles/tailscale/templates/etc/systemd/system/tailscaled.service.d/override.conf.j2
# Tailscale Systemd Service Override Template
# Consolidated override configuration for daemon config and custom settings
# Reference: https://tailscale.com/kb/1654/tailscaled-config-file
{% set unit_directives = ['After', 'Before', 'Wants', 'Requires', 'PartOf', 'Conflicts', 'ReloadPropagatedFrom', 'BindsTo', 'Requisite', 'OnFailure', 'PropagatesReloadTo', 'JoinsNamespaceOf'] %}
{% set unit_settings = {} %}
{% set service_settings = {} %}
{% for key, value in tailscale_service_override.items() %}
{% if key in unit_directives %}
{% set _ = unit_settings.update({key: value}) %}
{% else %}
{% set _ = service_settings.update({key: value}) %}
{% endif %}
{% endfor %}

{% if unit_settings %}
[Unit]
{% for section, settings in unit_settings.items() %}
{% if settings is iterable and settings is not string %}
# {{ section }} list values
{% for item in settings %}
{{ section }}={{ item }}
{% endfor %}
{% else %}
{{ section }}={{ settings }}
{% endif %}
{% endfor %}

{% endif %}
[Service]
{% if service_settings %}
{% for section, settings in service_settings.items() %}
{% if section == 'Environment' and settings is mapping %}
# Environment variables
{% for key, value in settings.items() %}
Environment={{ key }}={{ value }}
{% endfor %}
{% elif section == 'ExecStartPre' and settings is iterable and settings is not string %}
# Pre-execution commands
{% for command in settings %}
ExecStartPre={{ command }}
{% endfor %}
{% elif settings is mapping %}
# {{ section }} section
{% for key, value in settings.items() %}
{% if key == 'Environment' and value is mapping %}
# Environment variables as key-value pairs
{% for env_key, env_value in value.items() %}
Environment={{ env_key }}={{ env_value }}
{% endfor %}
{% elif key == 'Environment' and value is iterable and value is not string %}
# Environment variables as list
{% for env_var in value %}
Environment={{ env_var }}
{% endfor %}
{% elif key == 'ExecStartPre' and value is iterable and value is not string %}
# Pre-execution commands
{% for command in value %}
ExecStartPre={{ command }}
{% endfor %}
{% else %}
{{ key }}={{ value }}
{% endif %}
{% endfor %}
{% elif settings is iterable and settings is not string %}
# {{ section }} list values
{% for item in settings %}
{{ section }}={{ item }}
{% endfor %}
{% else %}
{{ section }}={{ settings }}
{% endif %}
{% endfor %}
{% endif %}
